# rootfs/etc/s6-overlay/s6-rc.d/svc-duckdns/run
#!/command/with-contenv bash
set -euo pipefail

echo "[duckdns] Starting DuckDNS updater (every 60 seconds)"

while true; do
    /usr/local/bin/duckdns-update || true
    sleep 60
done
