# rootfs/etc/s6-overlay/s6-rc.d/svc-wireguard/finish
#!/command/with-contenv bash
echo "[wireguard] Stopping WireGuard interface wg0"
wg-quick down /config/wg/wg0.conf 2>/dev/null || true
